# ansible/roles/kiosk/files/browser-fullscreen.service
[Unit]
Description=Browser Fullscreen
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/xte -x $DISPLAY 'key F11'
