# ansible/roles/kiosk/templates/display_on.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Display On Service
PartOf=graphical-session.target
ConditionlArchitecture=arm

[Service]
Type=oneshot
ExecStart={{ kiosk_display_on_cmd }}
